# debian/prerm
#!/bin/sh
# Stop and un-enable the user service before the files disappear.
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl --global disable meeting-recorder.service >/dev/null 2>&1 || true
        # Stop any instance running for the invoking user (best effort).
        if [ -n "${SUDO_USER:-}" ]; then
            su - "$SUDO_USER" -c \
               'systemctl --user stop meeting-recorder.service' >/dev/null 2>&1 || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
