# packaging/meeting-recorder.service
[Unit]
Description=Smart Meeting Recorder (auto-detect meetings and offer to record)
Documentation=man:meeting-recorder(1)
After=graphical-session.target pipewire.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/meeting-recorder run
Environment=PYTHONUNBUFFERED=1
# Fallback X11 display for notifications and x11grab capture. The graphical
# session normally imports the real DISPLAY/XAUTHORITY -- and, on Wayland,
# WAYLAND_DISPLAY/XDG_SESSION_TYPE -- into the user manager, which is what
# selects the capture backend. This is only a last resort for X11 sessions.
Environment=DISPLAY=:0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
