# packaging/postinst
#!/bin/sh
# Enable the per-user background service for all users on this system.
set -e

case "$1" in
    configure)
        # --global enables the unit for every user session (it is a user unit,
        # started on graphical login). Never fail the install if systemd is not
        # available (containers, chroots).
        systemctl --global enable meeting-recorder.service >/dev/null 2>&1 || true

        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database -q /usr/share/applications || true
        fi
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -qtf /usr/share/icons/hicolor || true
        fi

        echo "meeting-recorder: starts automatically at your next login."
        echo "  Start now:  systemctl --user start meeting-recorder"
        echo "  Settings:   meeting-recorder settings"
        ;;
esac

exit 0
